# container-images/git-gc/cron/crontab
0 6 * * * /var/tools/gc-all.sh
30 17 * * * /var/tools/gc-all.sh
